# docs/yolo-inference.service
[Unit]
Description=YOLO Inference Pipeline (Jetson / yolo-saad container)
Documentation=https://github.com/requiem002/esa-edge-perception
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=esa_jetson
Environment=HOME=/home/esa_jetson

# The yolo-saad container was committed from yolo-orin:v2.
# We do NOT recreate it here — only start it if stopped.
ExecStartPre=/usr/bin/docker start yolo-saad

# Headless inference: saves output/detect_*.avi and detections_*.csv
# under ~/Desktop/output (which is /workspace/output inside the container).
ExecStart=/usr/bin/docker exec yolo-saad python3 /workspace/inference.py --no-display

# Graceful stop, then SIGKILL after timeout.
ExecStop=/usr/bin/docker exec yolo-saad pkill -TERM -f "/workspace/inference.py"
KillMode=mixed
TimeoutStopSec=15

# Auto-restart on crash with backoff to avoid tight loops.
Restart=on-failure
RestartSec=10

# Logs land in journald — view with: journalctl -u yolo-inference.service -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yolo-inference

[Install]
WantedBy=multi-user.target
